# docker/supervisord.conf
; KusumaVision NMS — proses di container app (di-include /etc/supervisor/conf.d/).
; Menyalakan web (php-fpm + nginx) dan 3 daemon yang sama seperti install.sh:
; queue worker, scheduler, telnet proxy. Semua log diarahkan ke stdout container
; supaya terlihat di `docker compose logs app`.

[program:php-fpm]
command=php-fpm -F
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:kusumavision-worker]
command=php /var/www/html/artisan queue:work redis --tries=1 --max-time=3600
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
priority=20
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:kusumavision-scheduler]
command=php /var/www/html/artisan schedule:work
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:kusumavision-telnet-proxy]
command=php /var/www/html/artisan telnet:proxy
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
priority=20
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
